# petalinux/project-spec/meta-user/recipes-apps/webfwu/busybox/busybox-httpd.service.in
[Unit]
Description=Busybox HTTP daemon
After=network.target

[Service]
Type=forking
ExecStart=/usr/sbin/httpd -h /srv/www
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target